# deploy/systemd/blog-deploy.path
[Unit]
Description=Watch for an authenticated blog deployment trigger

[Path]
PathChanged=/srv/blog/.deploy-trigger
Unit=blog-deploy.service

[Install]
WantedBy=multi-user.target
